# node/systemd/services/solana-healthcheck.service
[Unit]
Description=Solana RPC healthcheck

[Service]
Type=oneshot
ExecStart=/usr/local/bin/solana-healthcheck.sh
