# node/systemd/services/solana-ledger-cleanup.service
[Unit]
Description=Solana ledger snapshot/remote cleanup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/solana-ledger-cleanup.sh
