# node/scripts/solana-mev-simple.service
[Unit]
Description=Solana MEV System
After=network.target

[Service]
Type=simple
User=kidgordones
WorkingDirectory=/home/kidgordones/0solana/node
ExecStart=/bin/bash /home/kidgordones/0solana/node/scripts/start-services.sh
ExecStop=/bin/bash -c 'lsof -ti:3001,8055,8081,8083,8085 | xargs -r kill -9 2>/dev/null || true'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target