# systemd/backend-api.service
[Unit]
Description=Solana MEV Backend API (FastAPI/Uvicorn)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=solana2
WorkingDirectory=/home/kidgordones/0solana/solana2/api
EnvironmentFile=/home/kidgordones/0solana/solana2/configs/backend.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8080 --workers 4 --http httptools --loop uvloop --log-level info
Restart=always
RestartSec=5
TimeoutStartSec=30
KillSignal=SIGTERM
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target