# scripts/solana-mev.service
[Unit]
Description=Solana MEV System - Ultra High Performance Trading Infrastructure
Documentation=https://github.com/kidgordones/0solana
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=forking
User=kidgordones
Group=kidgordones
WorkingDirectory=/home/kidgordones/0solana/node

# Environment variables
Environment="NODE_ENV=production"
Environment="RUST_LOG=info"
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/home/kidgordones/.cargo/bin:/home/kidgordones/.local/bin"

# Pre-start cleanup
ExecStartPre=/bin/bash -c 'lsof -ti:3001,8055,8081,8083,8085 | xargs -r kill -9 2>/dev/null || true'
ExecStartPre=/bin/sleep 2

# Start all services
ExecStart=/home/kidgordones/0solana/node/scripts/start-all-services.sh

# Stop command
ExecStop=/bin/bash -c 'kill $(cat /tmp/*.pid 2>/dev/null) 2>/dev/null || true'
ExecStop=/bin/bash -c 'lsof -ti:3001,8055,8081,8083,8085 | xargs -r kill -9 2>/dev/null || true'

# Restart policy
Restart=always
RestartSec=10
TimeoutStartSec=300

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

# Performance tuning
CPUWeight=100
IOWeight=100
TasksMax=4096

# Security
PrivateTmp=false
ProtectSystem=false
ProtectHome=false
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target